# scripts/launchd/io.github.gwanggyo-barun.clinic-content.daily-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent: silent daily git fetch for clinic-content-system.

  Installed by scripts/bootstrap-machine.sh to ~/Library/LaunchAgents/.
  The bootstrap script substitutes @@HOME@@ and @@REPO_DIR@@ before
  copying so the file is machine-specific.

  Schedule: weekdays 09:00 local time (Mon=1 through Fri=5).
  Behavior on a sleeping mac: launchd defers until the next wake, so
  the fetch will run when the machine resumes if it was off at 09:00.

  ──────────────────────────────────────────────────────────────────
  ⚠️ OPERATOR NOTES (added 2026-06-10 — see docs/launchd-notion-audit-20260610.md)

  • THIS IS A TEMPLATE, NOT A READY plist. It still contains the
    @@HOME@@ / @@REPO_DIR@@ placeholders. Install ONLY via
    `bash scripts/bootstrap-machine.sh`, which sed-substitutes them.
    Do NOT `launchctl bootstrap` this file directly — launchd would
    fail to resolve the paths.

  • DO NOT hand-copy this plist into ~/Library/LaunchAgents to "register
    a second copy". bootstrap-machine.sh already prevents duplicate
    registration (it does `launchctl bootout` before `bootstrap`), so
    re-running bootstrap is idempotent and self-replacing.

  • SCOPE: this job runs a READ-ONLY `git fetch` only (scripts/daily-fetch.sh).
    It never writes files, never calls the Notion API, and therefore is
    NOT a cause of any Notion orphan/duplicate-row issue. Look at the
    build -> Notion sync pipeline (build.py / shared/_notion_sync.py /
    .github/workflows) for those, per the audit report above.
  ──────────────────────────────────────────────────────────────────
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.gwanggyo-barun.clinic-content.daily-fetch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@REPO_DIR@@/scripts/daily-fetch.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@@REPO_DIR@@</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>@@HOME@@/Library/Logs/clinic-content-system/daily-fetch.log</string>

  <key>StandardErrorPath</key>
  <string>@@HOME@@/Library/Logs/clinic-content-system/daily-fetch.err</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
