# scripts/launchd/io.github.gwanggyo-barun.clinic-content.daily-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent: silent daily git fetch for clinic-content-system.

  Installed by scripts/bootstrap-machine.sh to ~/Library/LaunchAgents/.
  The bootstrap script substitutes @@HOME@@ and @@REPO_DIR@@ before
  copying so the file is machine-specific.

  Schedule: weekdays 09:00 local time (Mon=1 through Fri=5).
  Behavior on a sleeping mac: launchd defers until the next wake, so
  the fetch will run when the machine resumes if it was off at 09:00.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.gwanggyo-barun.clinic-content.daily-fetch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@REPO_DIR@@/scripts/daily-fetch.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@@REPO_DIR@@</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>@@HOME@@/Library/Logs/clinic-content-system/daily-fetch.log</string>

  <key>StandardErrorPath</key>
  <string>@@HOME@@/Library/Logs/clinic-content-system/daily-fetch.err</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
